# packages/falter-common/files/etc/uci-defaults/freifunk-berlin-ffwizard3-json
#!/bin/sh

# . /lib/functions/guard.sh
# guard "ffwizard3"

# check for translation-script
if ! which uci2ffwizard; then
    exit 1
fi

# retry on next boot if the wizard hasn't been completed yet
if [ ! -f /etc/config/ffwizard ]; then
    exit 1
fi

# create a ffwizard3-file, if there is none already
if [ ! -f /etc/ffwizard3.json ]; then
    uci2ffwizard >/etc/ffwizard3.json
    exit 0
fi
