# packages/falter-common/files/etc/uci-defaults/freifunk-berlin-network-ip-route
#!/bin/sh

# shellcheck disable=SC1091

. /lib/functions/guard.sh
guard "iproute"

echo 120 ffuplink >>/etc/iproute2/rt_tables
