# packages/falter-common/files/etc/profile.d/20_aliases.sh
#!/bin/sh

alias la="ls -lha"
alias neigh.sh="olsrd-neigh"
