# packages/falter-common/files/etc/rc.local
#!/bin/sh /etc/rc.common

# shellcheck disable=SC1090

[ -d /etc/rc.local.d ] && {
  for file in /etc/rc.local.d/*; do
    test -f "$file" && . "$file"
  done
}

exit 0;
