# packages/falter-policyrouting/files/etc/uci-defaults/freifunk-policyrouting
#!/bin/sh

# shellcheck source=/dev/null
. /lib/functions/guard.sh
guard "freifunk-policyrouting"

uci batch <<-EOF
        add ucitrack freifunk-policyrouting
        add_list ucitrack.@freifunk-policyrouting[-1].exec="/etc/init.d/freifunk-policyrouting restart"
        commit ucitrack
EOF

