# packages/falter-common/files/etc/init.d/freifunk
#!/bin/sh /etc/rc.common
# Freifunk Init
# $Id$

# START var is mandatory part of openwrt init-script syntax, though it appears unused
# shellcheck disable=SC2034
# don't follow that file in line 21...
# shellcheck disable=SC1090

START=99

boot() {
	test -f /etc/crontabs/root || touch /etc/crontabs/root

	grep -q 'killall -HUP dnsmasq' /etc/crontabs/root || {
		echo "*/5 * * * *	killall -HUP dnsmasq" >> /etc/crontabs/root
	}

	[ -d /etc/rc.local.d ] && {
		for file in /etc/rc.local.d/*; do
			test -f "$file" && . "$file"
		done
	}

	/etc/init.d/cron restart &
}
